# payApp.conf
[program:payApp]
command=/usr/local/bin/gunicorn wsgi -c gunicorn.config.py
directory=/var/www/payApp
autostart=true
autorestart=true